# scripts/nba-secrets-monitor.service
[Unit]
Description=NBA MCP Secrets Health Monitor
After=network.target

[Service]
Type=simple
User=ryanranft
Group=staff
WorkingDirectory=/Users/ryanranft/nba-mcp-synthesis
ExecStart=/usr/bin/python3 /Users/ryanranft/nba-mcp-synthesis/mcp_server/secrets_health_monitor.py --interval 300
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/Users/ryanranft/nba-mcp-synthesis
Environment=NBA_MCP_CONTEXT=WORKFLOW

[Install]
WantedBy=multi-user.target

